# build/crontab
0 0 * * SUN /root/update-roothints.sh > /var/log/update-roothints.log
30 0 * * SUN /root/update-applications.sh > /var/log/update-applications.log
